# dotfiles/.config/systemd/user/app-com.github.wwmm.easyeffects@autostart.service.d/delay.conf
[Service]
ExecStartPre=/usr/bin/sleep 15
